# s6/dbus/run
#!/bin/sh

exec dbus-daemon --system --nofork --nosyslog
