# s6/nmbd/run
#!/bin/sh

exec /usr/sbin/nmbd --foreground --no-process-group --log-stdout --debuglevel=1
