# s6/.s6-svscan/finish
#!/bin/sh

# do nothing
echo "Executing .s6-svscan/finish with arguments ${*}"
