# s6/smbd/run
#!/bin/sh

DEBUG_LEVEL="${DEBUG_LEVEL:-1}"

exec /usr/sbin/smbd --foreground --no-process-group --debug-stdout --debuglevel="${DEBUG_LEVEL}" < /dev/null
