# s6/smbd/run
#!/bin/sh

exec /usr/sbin/smbd --foreground --no-process-group --log-stdout --debuglevel=1
